# dbbackup/crontab
0 0 * * * /backup.sh>> /var/log/backup.log

